# rootfs/etc/s6-overlay/s6-rc.d/bot/run
#!/command/execlineb -P

with-contenv

fdmove -c 2 1
s6-setuidgid botuser
python3 -m hmnd_bot
